# tests/supervisord_test001.conf
[inet_http_server]
port=:9074

[supervisord]
nodaemon=true

[supervisorctl]
serverurl=http://localhost:9074

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[rpcinterface:multivisor]
supervisor.rpcinterface_factory = multivisor.rpc:make_rpc_interface
bind=*:9073

[group:Vacuum]
programs:vacuum_OH1,vacuum_OH2_1,vacuum_OH2_2

[group:PLC]
programs:wcid00a,wcid00b,wcid00c,wcid00d

[group:Counter]
programs:exits_2s,exits_10s,exit_1s_restart

[program:vacuum_OH1]
command=examples/full_example/vacuum OH1
autorestart=unexpected
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0
stderr_logfile=%(here)s/log/%(program_name)s.err.log
stderr_logfile_maxbytes=128kB
stderr_logfile_backups=0
startsecs=5

[program:vacuum_OH2_1]
command=examples/full_example/vacuum OH2_1
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0
startsecs=10

[program:vacuum_OH2_2]
command=examples/full_example/vacuum OH2_2
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0

[program:wcid00a]
command=examples/full_example/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0

[program:wcid00b]
command=examples/full_example/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0

[program:wcid00c]
command=examples/full_example/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0

[program:wcid00d]
command=examples/full_example/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0

[program:exits_2s]
command=examples/full_example/exits 2
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0
startsecs=5

[program:exits_10s]
command=examples/full_example/exits 10
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0
startsecs=1

[program:exit_1s_restart]
command=examples/full_example/exits 1
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=128kB
stdout_logfile_backups=0
startsecs=2
