# docker/supervisord/baslid001.conf
[inet_http_server]
port=:9001

[supervisord]
identifier=baslid001
nodaemon = true
pidfile = /var/run/supervisord.pid
logfile = /var/log/supervisord/supervisord.log
childlogdir = /var/log/supervisord
logfile_backups=10
logfile_maxbytes=1MB

[supervisorctl]
serverurl=http://localhost:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[eventlistener:multivisor-rpc]
command=multivisor-rpc
events=PROCESS_STATE,SUPERVISOR_STATE_CHANGE

[group:BeamViewer]
programs=basler_mbv1,basler_mbv2,basler_mbv3,basler_wbv1,basler_wbv2

[program:basler_wbv1]
command=Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:basler_wbv2]
command=Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:basler_mbv1]
command=Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:basler_mbv2]
command=Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:basler_mbv3]
process_name=basler_mbv3-%(process_num)s
numprocs=3
command=Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:talks_too_much]
command=talkative
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
