# docker/supervisord/lid002.conf
[inet_http_server]
port=:9001

[supervisord]
identifier=lid001
nodaemon = true
pidfile = /var/run/supervisord.pid
logfile = /var/log/supervisord/supervisord.log
childlogdir = /var/log/supervisord
logfile_backups=10
logfile_maxbytes=1MB

[supervisorctl]
serverurl=http://localhost:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[eventlistener:multivisor-rpc]
command=multivisor-rpc
events=PROCESS_STATE,SUPERVISOR_STATE_CHANGE

[group:Vacuum]
programs:vacuum_EH

[group:PLC]
programs:wcid00e,wcid00f,wcid00g

[group:Counter]
programs:P201

[program:vacuum_EH]
command=vacuum EH
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00e]
command=wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00f]
command=wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00g]
command=wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:P201]
command=CT2 EH
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
