# docker/supervisord/lid001.conf
[inet_http_server]
port=:9001

[supervisord]
identifier=lid001
nodaemon = true
pidfile = /var/run/supervisord.pid
logfile = /var/log/supervisord/supervisord.log
childlogdir = /var/log/supervisord
logfile_backups=10
logfile_maxbytes=1MB

[supervisorctl]
serverurl=http://localhost:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[eventlistener:multivisor-rpc]
command=multivisor-rpc
events=PROCESS_STATE,SUPERVISOR_STATE_CHANGE

[group:Vacuum]
programs:vacuum_OH1,vacuum_OH2_1,vacuum_OH2_2

[group:PLC]
programs:wcid00a,wcid00b,wcid00c,wcid00d

[group:Counter]
programs:exits_2s,exits_10s,exit_1s_restart

[program:multivisor]
command=multivisor -c /etc/multivisor.conf
autorestart=unexpected
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=%(program_name)s.err.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1

[program:vacuum_OH1]
command=vacuum OH1
autorestart=unexpected
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=%(program_name)s.err.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1
startsecs=5

[program:vacuum_OH2_1]
command=vacuum OH2_1
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
startsecs=10

[program:vacuum_OH2_2]
command=vacuum OH2_2
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00a]
command=wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00b]
command=wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00c]
command=wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00d]
command=wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:exits_2s]
command=exits 2
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
startsecs=5

[program:exits_10s]
command=exits 10
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
startsecs=1

[program:exit_1s_restart]
command=exits 1
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
startsecs=2
