# examples/full_example/supervisord_lid001.conf
[inet_http_server]
port=:9011

[supervisord]
logfile_backups=10
logfile_maxbytes=1MB
logfile=%(here)s/log/supervisor_lid001.log
pidfile=%(here)s/supervisor_lid001.pid
childlogdir=%(here)s/log
identifier=lid001

[supervisorctl]
serverurl=http://localhost:9011

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

#[rpcinterface:multivisor]
#supervisor.rpcinterface_factory = multivisor:make_rpcinterface
#name="Tiago's multivisor"
#port=:22000
#node.lid002.url = http://lid002:9002

[group:Vacuum]
programs:vacuum_OH1,vacuum_OH2_1,vacuum_OH2_2

[group:PLC]
programs:wcid00a,wcid00b,wcid00c,wcid00d

[group:Counter]
programs:exits_2s,exits_10s,exit_1s_restart

[program:vacuum_OH1]
command=%(here)s/vacuum OH1
autorestart=unexpected
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=%(here)s/log/%(program_name)s.err.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1
startsecs=5

[program:vacuum_OH2_1]
command=%(here)s/vacuum OH2_1
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
startsecs=10

[program:vacuum_OH2_2]
command=%(here)s/vacuum OH2_2
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00a]
command=%(here)s/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00b]
command=%(here)s/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00c]
command=%(here)s/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00d]
command=%(here)s/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:exits_2s]
command=%(here)s/exits 2
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
startsecs=5

[program:exits_10s]
command=%(here)s/exits 10
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
startsecs=1

[program:exit_1s_restart]
command=%(here)s/exits 1
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
startsecs=2
