# examples/full_example/supervisord_baslid001.conf
[inet_http_server]
port=:9031

[supervisord]
logfile_backups=10
logfile_maxbytes=1MB
logfile=%(here)s/log/supervisor_baslid001.log
pidfile=%(here)s/supervisor_baslid001.pid
childlogdir=%(here)s/log
identifier=baslid001

[supervisorctl]
serverurl=http://localhost:9031

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[rpcinterface:multivisor]
supervisor.rpcinterface_factory = multivisor.rpc:make_rpc_interface
bind=*:9032

[group:BeamViewer]
programs=basler_mbv1,basler_mbv2,basler_mbv3,basler_wbv1,basler_wbv2

[program:basler_wbv1]
command=%(here)s/Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:basler_wbv2]
command=%(here)s/Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:basler_mbv1]
command=%(here)s/Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:basler_mbv2]
command=%(here)s/Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:basler_mbv3]
process_name=basler_mbv3-%(process_num)s
numprocs=3
command=%(here)s/Lima %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:talks_too_much]
command=%(here)s/talkative
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
