# examples/full_example/supervisord_lid002.conf
[inet_http_server]
port=:9021

[supervisord]
logfile_backups=10
logfile_maxbytes=1MB
logfile=%(here)s/log/supervisor_lid002.log
pidfile=%(here)s/supervisor_lid002.pid
childlogdir=%(here)s/log
identifier=lid002

[supervisorctl]
serverurl=http://localhost:9021

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[rpcinterface:multivisor]
supervisor.rpcinterface_factory = multivisor.rpc:make_rpc_interface
bind=*:9022

[group:Vacuum]
programs:vacuum_EH

[group:PLC]
programs:wcid00e,wcid00f,wcid00g

[group:Counter]
programs:P201

[program:vacuum_EH]
command=%(here)s/vacuum EH
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/vacuum_EH.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00e]
command=%(here)s/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00f]
command=%(here)s/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:wcid00g]
command=%(here)s/wago %(program_name)s
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:P201]
command=%(here)s/CT2 EH
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(here)s/log/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
